# mihomo
#!/sbin/openrc-run
# Copyright 1999-2018 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

CONFIG="/etc/mihomo/config.yaml"
PIDFILE="/run/mihomo.pid"
LOGFILE="/var/log/mihomo.log"

depend() {
	need net
}

checkconfig() {
	if [ ! -f ${CONFIG} ]; then
		ewarn "${CONFIG} does not exist."
	fi
}

start() {
	checkconfig || return 1

	ebegin "Starting mihomo"
	ebegin "Log File : ${LOGFILE}"
	start-stop-daemon --start	\
	-b -1 ${LOGFILE} -2 ${LOGFILE}	\
	-m -p ${PIDFILE}		\
	--exec /usr/bin/mihomo -f ${CONFIG}
	eend $?
}

stop() {
	ebegin "Stopping mihomo"
	start-stop-daemon --stop -p ${PIDFILE}
	eend $?
}
